# scripts/blog-platform.service
[Unit]
Description=Blog Platform Web Application
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/blog
ExecStart=/usr/bin/bash /var/www/blog/scripts/start-dev.sh
Restart=always
RestartSec=10
StartLimitBurst=5
StartLimitIntervalSec=60

# 环境变量
Environment=NODE_ENV=production
Environment=PORT=3000

# 日志
StandardOutput=append:/var/www/blog/logs/web/service.log
StandardError=append:/var/www/blog/logs/web/error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
